# systemd/bootupd.socket
[Socket]
ListenStream=/run/bootupd.sock
SocketMode=0600

[Install]
WantedBy=sockets.target
